# services/crt-web.service
[Unit]
Description=CRT Kitchen TV Web Config UI
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/opt/crt-kitchen-tv
User=crt
Group=crt
Environment=CRT_CONFIG=/etc/crt-kitchen-tv/config.yaml
ExecStart=/opt/crt-kitchen-tv/venv/bin/python -m waitress --listen=0.0.0.0:8080 server.app:app
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
