# services/crt-ui.service
[Unit]
Description=CRT Kitchen TV Frontend UI (X11)
After=network.target
Wants=network.target
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
User=crt
Group=crt
WorkingDirectory=/opt/crt-kitchen-tv

# Avoid audio init blocking UI bring-up; configure ALSA/ReSpeaker later.
Environment=SDL_AUDIODRIVER=dummy

# Attach to tty1 and let startx run the session using /home/crt/.xinitrc
TTYPath=/dev/tty1
StandardInput=tty-force
StandardOutput=journal
StandardError=journal

ExecStart=/usr/bin/startx -- vt1 -nocursor
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
